# dynamodb/dynamodb.conf
[program:dynamodb]
command=/usr/bin/start-dynamodb.sh
autostart=true
autorestart=true
