# zookeeper/zookeeper.conf
[program:zookeeper]
command=/usr/bin/start-zk.sh
autostart=true
autorestart=true
