# hadoop/hadoop.conf
[program:hadoop]
command=/usr/bin/start-hadoop.sh -d
autostart=true
autorestart=true
