# hadoop/hadoop.conf
[program:hadoop]
command=/usr/bin/start-hadoop.sh
autostart=true
autorestart=true
