# hive/hive.conf
[program:hive-metadata]
command=/usr/bin/start-metadata.sh
autostart=true
autorestart=true
priority=1

[program:hive-server]
command=/usr/bin/start-hive.sh
autostart=true
autorestart=true
priority=2
